# .config/sv/picom/run
#!/bin/sh
exec 2>&1
exec picom --experimental-backends --config ~/.config/picom/picom.conf
