# .config/sv/runsvdir-yusuf/finish
#!/bin/sh

sv -w600 force-stop /home/yusuf/service/*
sv exit /home/yusuf/service/*
