# .config/sv/runsvdir-yusuf/run
#!/bin/sh

[ ! -d /run/runit.yusuf ] && \
  mkdir -p /run/runit.yusuf && \
  chown duncan:duncan /run/runit.yusuf

export HOME=/home/yusuf
groups=$(id -Gn yusuf | tr ' ' ':')

exec 2>&1
exec chpst -u "yusuf:$groups" sh -l -c "exec runsvdir /home/yusuf/.config/service 'log: ...........................................................................................................................................................................................................................................................................................................................................................................................................'" yusuf
